# .config/autostart/libinput-gestures.desktop
[Desktop Entry]
Type=Application
NoDisplay=true
Terminal=false
Name=Libinput Gestures
Exec=/usr/bin/libinput-gestures
Icon=libinput-gestures
Comment=Background application to intercept and action libinput gestures from touchpad.
Categories=System;
